# 90-usbmidiconnect.rules
ACTION=="add",\
SUBSYSTEMS=="usb",\
RUN+="/bin/bash /usr/local/share/timidity/midiconnect.sh"